# docker/django/celery.crontab
30 6 * * * cd /data/web && /data/web/docker/django/run_management_command.sh taskstore update_statistics > /tmp/update_stats.log 2>&1
30 7 * * * cd /data/web && /data/web/docker/django/run_management_command.sh taskstore gc_large_repos > /tmp/gc_large_repos.log 2>&1
30 11 * * * cd /data/web && /data/web/docker/django/run_management_command.sh taskstore refresh_certificates > /tmp/refresh_certificates.log 2>&1
30 12 * * * cd /data/web && /data/web/docker/django/run_management_command.sh taskstore delete_old_accounts > /tmp/delete_old_accounts.log 2>&1
30 13 * * * cd /data/web && /data/web/docker/django/run_management_command.sh taskstore delete_old_changes > /tmp/delete_old_changes.log 2>&1
