# docker/django/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null

[program:django]
command=/data/web/run.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
stopsignal=QUIT
