# scripts/vagrant/simple_celery_upstart.conf
description "Inthe.AM Celery"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 10 5
console log

script
    exec su -c /var/www/twweb/start_celery.sh vagrant
end script
